# roles/monitoring_agent/templates/systemd/nginx-exporter/nginx_exporter.socket.j2
[Unit]
Description=Nginx Prometheus Exporter

[Socket]
ListenStream={{ monitoring_agent_nginx_exporter_port }}

[Install]
WantedBy=sockets.target
