# roles/monitoring_agent/templates/systemd/node-exporter/node_exporter.socket.j2
[Unit]
Description=Node Exporter

[Socket]
ListenStream={{ monitoring_agent_node_exporter_port }}

[Install]
WantedBy=sockets.target
